# supervisord.conf
[supervisord]
nodaemon=true

[program:ventrilo]
user=ventrilo
directory=/srv/ventrilo/ventsrv/
command=bash -c "/srv/ventrilo/ventsrv/ventrilo_srv -d"
autostart=true